# deploy/blockshelf.service
[Unit]
Description=Gunicorn for BlockShelf (Django)
After=network.target

[Service]
# <-- change this to your actual user that owns ~/BlockShelf
User=%i
Group=%i

# Repo location: /home/<user>/BlockShelf
WorkingDirectory=/home/%i/BlockShelf
Environment="PYTHONUNBUFFERED=1"
# Optional: use a .env in your repo root
EnvironmentFile=/home/%i/BlockShelf/.env

# Virtualenv inside the repo
ExecStart=/home/%i/BlockShelf/.venv/bin/gunicorn \
  -c /home/%i/BlockShelf/gunicorn.conf.py \
  blockshelf_inventory.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Create /run/blockshelf for the socket
RuntimeDirectory=blockshelf
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
